# projects/Generic/filesystem/usr/lib/systemd/system/pulse8-cec-inputattach@.service
[Unit]
Description=inputattach for pulse8-cec device on %I

[Service]
Type=simple
ExecStart=/usr/bin/inputattach --pulse8-cec /dev/%I
