# projects/Generic/filesystem/usr/lib/systemd/system/rainshadow-cec-inputattach@.service
[Unit]
Description=inputattach for rainshadow-cec device on %I

[Service]
Type=simple
ExecStart=/usr/bin/inputattach --rainshadow-cec /dev/%I
