# packages/addons/service/nfs-ganesha/source/system.d/service.nfs-ganesha.service
[Unit]
Description=nfs-ganesha
After=network-online.target
Requires=network-online.target

[Service]
Type=forking
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.nfs-ganesha/bin/nfs-ganesha.start"
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
