# packages/addons/service/zabbix-proxy/source/system.d/service.zabbix-proxy.service
[Unit]
Description=Zabbix Proxy service
After=network-online.target
Requires=network-online.target

[Service]
Type=notification
ExecStart=/bin/sh -c "exec sh /storage/.kodi/addons/service.zabbix-proxy/bin/zabbix-proxy.start"
ExecStop=/bin/sh -c "exec sh /storage/.kodi/addons/service.zabbix-proxy/bin/zabbix-proxy.stop"
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/zabbix/zabbix_proxy.pid

[Install]
WantedBy=kodi.target
